# etc/supervisord.d/flink.ini
[program:flink]
command = /opt/qnib/flink/bin/start_flink.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
